# recipes-ti/beagleboard/gadget-init/bone-gmass-eject.rules
SUBSYSTEM=="block",ACTION=="change",KERNEL=="mmcblk0p1",RUN+="/usr/sbin/g-ether-insert.sh"
